# roles/svgcircus_api/files/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[supervisorctl]
serverurl=http://127.0.0.1:8998


[program:gunicorn]
command=/opt/svgcircus_api/venv/bin/gunicorn svgcircus.wsgi:application --bind "127.0.0.1:8999"
directory=/opt/svgcircus_api
autorestart=true
redirect_stderr=true


